# alpineproc/etc/services.d/01-crontab/run
#!/bin/sh
echo "Start running crontab service ... ..." >&2
exec /usr/sbin/crond -f
echo "Finish running crontab service ... ..." >&2
